# templates/mcce-digest.service.j2
[Unit]
Description=MCCE Serial Digest Service
Wants=network-online.target
After=network.target network-online.target
OnFailure={{ systemd_mailer }}

[Service]
Type=oneshot
User=mcce
Group=mcce
WorkingDirectory={{ working_dir }}
Environment="MIN_TIME_RES={{ time_seconds }}"
Environment="DATABASE_URL={{ database_url }}"
ExecStart={{ python_path }} day_digest.py

[Install]
WantedBy=multi-user.target
