# templates/mcce-digest.timer.j2
[Unit]
Description=MCCE Serial Digest Timer

[Timer]
OnCalendar={{ schedule_time }}

[Install]
WantedBy=timers.target
